# deploy/wholesale.service
[Unit]
Description=Cambodia Wholesale Backend (Go)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wholesale/backend
ExecStart=/opt/wholesale/backend/wholesale
Restart=always
RestartSec=5
Environment=PATH=/usr/bin:/bin
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
